# HOME/Library/LaunchAgents/com.tmux.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.tmux.update</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/bash</string>
			<string>-c</string>
			<string>
				for d in $HOME/.tmux/plugins/*; do
					/run/current-system/sw/bin/git -C $d pull --recurse-submodules
				done
			</string>
		</array>
		<key>KeepAlive</key>
		<dict>
			<key>SuccessfulExit</key>
			<false/>
		</dict>
		<key>StandardErrorPath</key>
		<string>/tmp/tmuxupdateError.log</string>
		<key>StandardOutputPath</key>
		<string>/tmp/tmuxupdate.log</string>
	</dict>
</plist>
