# HOME/Library/LaunchAgents/com.nvim.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.nvim.update</string>
		<key>ProgramArguments</key>
		<array>
			<string>/run/current-system/sw/bin/nvim</string>
			<string>--headless</string>
			<string>+Lazy! sync</string>
			<string>+MasonUpdate</string>
			<string>+MasonToolsUpdateSync</string>
			<string>+qa</string>
		</array>
		<key>KeepAlive</key>
		<dict>
			<key>SuccessfulExit</key>
			<false/>
		</dict>
		<key>StandardErrorPath</key>
		<string>/opt/homebrew/var/log/nvimupdateError.log</string>
		<key>StandardOutputPath</key>
		<string>/opt/homebrew/var/log/nvimupdate.log</string>
	</dict>
</plist>
